# supervisord.conf
[supervisord]
nodaemon=true

[program:app]
command=python app.py

[program:airflow-standalone]
command=airflow standalone
